# package_tree/DEBIAN/postinst
#!/bin/bash

set -u
set -e


systemctl enable webilastik@-opt-webilastik.service
systemctl daemon-reload
systemctl restart webilastik@-opt-webilastik.service
systemctl restart nginx.service