# package_tree/DEBIAN/prerm
#!/bin/bash

set -u
set -e

systemctl stop webilastik@-opt-webilastik.service
systemctl disable webilastik@-opt-webilastik.service